# .config/systemd/user/backup.timer
[Unit]
Description=Schedule backups 

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
