# .config/systemd/user/backup.timer
[Unit]
Description=Schedule backups 

[Timer]
OnCalendar=*-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
