# .config/systemd/user/backup.service
[Unit]
Description=Executes a backup of the hard drive to backblaze
 
[Service]
Type=oneshot
ExecStart=~/.scripts/backup_backblaze.sh
