# .config/systemd/user/sway.service
[Unit]
Description=sway
[Service]
ExecStart=sway 2>/tmp/sway.log
[Install]
WantedBy=default.target
